# systemd/hl-recorder.service
[Unit]
Description=HL Reaper L2/OI microstructure recorder (mainnet, read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=veilcerpin
WorkingDirectory=/opt/hl-reaper
ExecStart=/opt/hl-reaper/venv/bin/python /opt/hl-reaper/scripts/run_recorder.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
