# systemd/hl-spot-poller.service
[Unit]
Description=HL Reaper Binance spot price poller/recorder (spot-perp lead/lag, read-only)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=veilcerpin
WorkingDirectory=/opt/hl-reaper
ExecStart=/opt/hl-reaper/venv/bin/python /opt/hl-reaper/scripts/run_spot_poller.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
